# ansible-playbook/roles/nas_server/templates/media-storage.mount.j2
[Unit]
Description=External hdd mount
After=local-fs.target

[Mount]
What={{external_hdd_dev}}
Where=/media/storage
Type=ext4
Options=defaults,rw,noexec

[Install]
WantedBy=default.target
